# v2/crontab
@daily cd ~/Documents/wetlands/v2 && ./run
